# utaskdaemon.service
[Unit]
Description=utaskdaemon

[Service]
ExecStart=/opt/utask/utaskdaemon

[Install]
WantedBy=default.target
